# systemd_units/finder.service
[Unit]
AssertPathExists=/home/foodfinder/opentelemetry-starter

[Service]
WorkingDirectory=/home/foodfinder/opentelemetry-starter
ExecStart=/home/foodfinder/opentelemetry-starter/run_finder.sh --is_prod
Restart=always
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=default.target
